# src/Apps/GardenControl/GardenControl.Server/Install/gardencontrol.service
[Unit]
Description=GardenControl Bewässerungssteuerung
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/gardencontrol
ExecStart=/home/pi/gardencontrol/GardenControl.Server
Restart=always
RestartSec=10
# Alle Ausgänge aus bei Crash
ExecStopPost=/usr/bin/bash -c 'echo 0 > /sys/class/gpio/gpio17/value 2>/dev/null; echo 0 > /sys/class/gpio/gpio22/value 2>/dev/null; echo 0 > /sys/class/gpio/gpio23/value 2>/dev/null; echo 0 > /sys/class/gpio/gpio27/value 2>/dev/null'

# Umgebung
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=DOTNET_RUNNING_IN_CONTAINER=false

[Install]
WantedBy=multi-user.target
